# jukebox.service
[Unit]
Description=Jukebox

[Service]
Type=simple
ExecStart=/srv/jukebox/start

[Install]
WantedBy=multi-user.target
